# service/vaderbot.conf
# vaderbot - run the vaderbot
#
# Check out the git repo and start the
# vaderbot.

description "Update and start vaderbot"

start on startup

exec /usr/local/bin/vaderbot.sh
